# infra/systemd/abs-purge-deleted-accounts.service
[Unit]
Description=ABS daily account purge (GDPR right-to-erasure)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/app
ExecStart=/usr/bin/python3 /app/infra/scripts/purge_deleted_accounts.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
